# udev/rules.d/90-keyboard-xset.rules
ACTION=="add", KERNEL=="event*", SUBSYSTEM=="input", RUN+="/etc/udev/scripts.d/keyboard-xset.sh"
