# udev/rules.d/90-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/flock -n /tmp/monitor.udev.lock /home/azat/dwm-assets/hotplug-monitor.sh"
